# scripts/daemons/agent-memory-proxy.service
[Unit]
Description=Agent Memory Proxy
After=network.target

[Service]
Type=simple
User=%USER%
Environment=%AGENT_MEMORY_PATHS%
ExecStart=/home/%USER%/.cargo/bin/uv run --directory %WORKING_DIR% amp
Restart=always
WorkingDirectory=%WORKING_DIR%

[Install]
WantedBy=multi-user.target
